# data/ubuntu-location-service.conf.in
description "Location Services"

start on started dbus

exec /usr/bin/ubuntu-location-serviced --bus system --provider gps::Provider
