# data/ubuntu-location-service.conf.in
description "Location Services"

start on started dbus

respawn

exec /usr/bin/ubuntu-location-serviced --bus system --provider gps::Provider --provider espoo::Provider  --espoo::Provider::name="com.ubuntu.espoo.Service.Provider" --espoo::Provider::path="/com/ubuntu/espooo/Service/Provider"
