# data/ubuntu-location-service.conf.in
description "Location Services"

start on started dbus

exec @libdir@/ubuntu-location-service --bus system --provider gps::Provider
