# data/ubuntu-location-service-trust-stored.conf.in
description "Location Services Trust Store Daemon"

start on started JOB=dbus and started JOB=unity8
stop on stopping JOB=dbus or (:sys:stopped JOB=ubuntu-location-service)

respawn

exec /usr/bin/trust-stored-skeleton \
    --remote-agent DBusRemoteAgent --bus=system \
    --local-agent MirAgent \
    --trusted-mir-socket=/var/run/user/$(id -u)/mir_socket_trusted \
    --for-service UbuntuLocationService \
    --store-bus session
