# data/ubuntu-location-service-trust-stored.conf.in
description "Location Services Trust Store Daemon"

start on started dbus

respawn

exec /usr/bin/trust-stored-skeleton \
    --remote-agent DBusRemoteAgent --bus=system \
    --local-agent MirAgent \
    --for-service @UBUNTU_LOCATION_SERVICE_TRUST_STORE_SERVICE_NAME@ \
    --store-bus session
